# systemd/.config/systemd/user/mail_work_inbox.service
[Unit]
Description=Check work inbox
Wants=mail_inbox.timer

[Service]
Type=simple
ExecStart=/home/markus/local/bin/mail_inbox.sh work

