# systemd/.config/systemd/user/testservice.service
[Unit]
Description=Prints the date every minute

[Service]
Type=simple
ExecStart=/bin/date

[Install]
WantedBy=mail_input.target
