# systemd/.config/systemd/user/xinitrc.target
[Unit]
Description=xinitrc stuff
Requires=default.target
After=default.target
